# fusion-sync/fusion-sync.service
[Unit]
Description=FusionSync daemon

[Service]
ExecStart=/usr/local/bin/fusion-sync daemon
Restart=always

[Install]
WantedBy=default.target
